# .config/autostart/fix-gesture-jump.desktop
[Desktop Entry]
Type=Application
Name=Fix Gesture Jump
Exec=~/.local/bin/fix_gesture_jump
StartupNotify=false
