# .config/autostart/workspace-monitor.desktop
[Desktop Entry]
Type=Application
Name=Workspace Monitor
Exec=~/.local/bin/workspace_monitor
StartupNotify=false